# supervisord.conf
[supervisord]
nodaemon=true

[program:nodeserv]
command=/usr/local/bin/node /src/app.js
stdout_logfile=/var/log/gs/server.log
redirect_stderr=true
logfile_maxbytes=50MB
